# roles/firefox_cache/files/browser-sync-periodic.service
[Unit]
Description=Run Periodic Sync for Browser Caches

[Service]
Type=oneshot
ExecStart=%h/.local/bin/browser-sync.sh sync
