# roles/firefox-cache/files/firefox-sync-periodic.service
[Unit]
Description=Run Periodic Sync for Firefox Cache

[Service]
Type=oneshot
ExecStart=%h/.local/bin/firefox-sync.sh sync
